# layers/meta-assistant/recipes-service/assistant/assistant-update.service
[Unit]
Description=Update assistant

[Service]
Restart=always
RestartSec=1800s
ExecStart=/bin/bash /home/user/script.sh

